# cron.d/gnome-shell-font
# Set GNOME Shell font to Ubuntu on each reboot
@reboot sed -i "s/cantarell/Ubuntu/gi" /usr/share/gnome-shell/theme/gnome-shell.css

